# adrsir-api.socket
[Unit]
Description=ADRSIR-API socket

[Socket]
ListenStream=/run/adrsir-api.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target